# wm/wallpaper/.config/systemd/user/wallpaper_timer.service
[Unit]
Description=Wallpaper changer

[Service]
ExecStart=/home/mori/.local/bin/wallpaper_timer.sh
