# wm/wallpaper/.config/systemd/user/wallpaper_timer.timer
[Unit]
Description=Runs wallpaper changer 10 minutes

[Timer]
OnCalendar=*:0/10
Persistent=true

[Install]
WantedBy=timers.target

